# crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
*/5 * * * * www-data bash /var/www/html/myddleware-cron.sh 2>&1
* * * * * www-data bash /var/www/html/script/monitoring.sh 2>&1
